# docker/supervisord.conf
[program:vnc]
user=student
environment=HOME="/home/student",USER="student"
command=/usr/local/bin/start-vnc.sh
directory=/home/student
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websockify]
user=student
environment=HOME="/home/student",USER="student"
command=/bin/bash -lc 'sleep 2 && exec websockify --web=/usr/share/novnc/ 0.0.0.0:6080 localhost:5901'
directory=/home/student
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:iot-device-core]
user=student
environment=HOME="/home/student",USER="student"
directory=/home/student/workspace/device_core
command=/bin/bash -lc 'sleep 6 && exec ./build/iot_device_core config.json'
autostart=true
autorestart=true
priority=25
startretries=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:flask-web]
user=student
environment=HOME="/home/student",USER="student",PYTHONUNBUFFERED="1",IOT_DB_PATH="/home/student/workspace/data/device.db"
directory=/home/student/workspace/web
command=/bin/bash -lc 'sleep 2 && exec python3 -m flask --app app run --host 0.0.0.0 --port 8080'
autostart=true
autorestart=true
priority=30
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
